# brln-rpc-server/brln-rpc-server.service
[Unit]
Description=BRLN-RPC-Server - Servidor HTTP Multi-Chain para BRLN-OS
Documentation=https://github.com/pagcoinbr/brln-os
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/brln-os/brln-rpc-server
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brln-rpc-server

# Environment
Environment=NODE_ENV=production
Environment=CONFIG_FILE=/root/brln-os/brln-rpc-server/config/config.json

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/brln-os/brln-rpc-server/logs
ReadWritePaths=/data/lnd

# Process limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
