# services/elementsd.service
[Unit]
Description=Elements Node Daemon
After=network.target

[Service]
User=admin
Group=admin
ExecStart=/usr/local/bin/elementsd -datadir=/data/elements
ExecStop=/usr/local/bin/elements-cli -datadir=/data/elements stop
Restart=on-failure
RestartSec=30
TimeoutStartSec=300
TimeoutStopSec=300
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
NoNewPrivileges=true
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
