# services/gotty-fullauto.service
[Unit]
Description=Terminal Web para BRLN FullAuto
After=network.target

[Service]
User=root
WorkingDirectory=/home/pagcoin/brln-os
Environment=TERM=xterm
ExecStart=/usr/local/bin/gotty -p 3131 -w bash /home/pagcoin/brln-os/brunel.sh
Restart=always

[Install]
WantedBy=multi-user.target
