# services/gotty-fullauto.service
[Unit]
Description=Terminal Web para BRLN FullAuto
After=network.target

[Service]
User=admin
WorkingDirectory=/home/admin
Environment=HOME=/home/admin
ExecStart=/usr/local/bin/gotty -w -p 3131 bash -c 'while true; do bash <(curl -s https://raw.githubusercontent.com/pagcoinbr/brlnfullauto/teste_v0.9/run.sh); done'
Restart=always

[Install]
WantedBy=multi-user.target
